# config-old/debian/cron.daily/privacyidea-rotate-audit
#!/bin/bash
python /usr/lib/python2.7/dist-packages/privacyidea/lib/auditmodules/sqlaudit.py -f /etc/privacyidea/privacyidea.ini

