# deploy/crontab/privacyidea
# Run the privacyidea-cron runner every 5 minutes
*/5 * * * *	privacyidea	privacyidea-cron run_scheduled -c

# Audit rotate 
#1 2 * * *	privacyidea	pi-manage audit rotate_audit --age 360
